# moneywiz/rootfs/etc/services.d/app/finish
#!/usr/bin/env bash
set -euxo pipefail
# vim: set ft=sh :
