# moneywiz/rootfs/etc/services.d/app/finish
#!/usr/bin/with-contenv sh
# vim: set ft=sh :
